# deploy/argusd.service
[Unit]
Description=ARGUS - Adaptive RDMA Guard & Utilization Sentinel
Documentation=https://github.com/KevinWeiss1995/ARGUS
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/argus/argusd.conf
ExecStart=/bin/sh -c 'exec /usr/local/bin/argusd ${ARGUS_CONFIG:+--config "$ARGUS_CONFIG"} --mode "$ARGUS_MODE" --ebpf-path "$ARGUS_EBPF_PATH" --metrics-addr "$ARGUS_METRICS_ADDR" --log-level "$ARGUS_LOG_LEVEL" --window-secs "$ARGUS_WINDOW_SECS" $ARGUS_EXTRA_ARGS'

Restart=on-failure
RestartSec=5

# Resource limits — ARGUS targets <1% CPU overhead
CPUQuota=5%
Nice=19
IOSchedulingClass=idle
MemoryMax=256M

# Security hardening — root is required for eBPF load, tracefs, kallsyms,
# and /proc/pid reads. security.rs drops all caps except CAP_BPF after init.
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_BPF CAP_PERFMON CAP_NET_ADMIN CAP_SYSLOG CAP_DAC_READ_SEARCH
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/argus /var/run/argus /sys/class/infiniband
ReadOnlyPaths=/sys/kernel /proc
PrivateTmp=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
SystemCallArchitectures=native

# Logging to journald (agent emits JSON when --tui is absent)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argusd

[Install]
WantedBy=multi-user.target
